# deploy/veriqrahq-api.service
[Unit]
Description=VeriqraHQ Compras Estatales API
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/veriqrahq/backend
EnvironmentFile=/var/www/veriqrahq/backend/.env
ExecStart=/var/www/veriqrahq/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8020
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
